# service/delta-bot-donchian.service
[Unit]
Description=Delta Exchange Trading Bot - Donchian Channel (PIUSD, PIPPINUSD, RIVERUSD, BERAUSD, PAXGUSD)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/delta-exchange-alog
# Runs all 5 Donchian coins in a single process using parallel threads.
# Each coin's logs still go to its own per-symbol file (logs/pi.log etc.)
# because _add_per_symbol_log_handler() in runner.py attaches a thread-filtered
# RotatingFileHandler per coin.  The shared DeltaRestClient serialises all API
# calls via the in-process threading.Lock in RateLimiter — no startup burst.
ExecStart=/home/pi/delta-exchange-alog/venv/bin/python3 -u /home/pi/delta-exchange-alog/run_terminal.py --strategy 14 --non-interactive
Restart=on-failure
RestartSec=60
EnvironmentFile=/home/pi/delta-exchange-alog/config/.env

# This captures Python-level crashes / uncaught exceptions from the main thread.
# Per-coin trade/signal logs are written directly to their own files by the
# thread-filtered log handlers (logs/pi.log, logs/pippin.log, etc.).
StandardOutput=append:/home/pi/delta-exchange-alog/logs/donchian-service.log
StandardError=append:/home/pi/delta-exchange-alog/logs/donchian-service.error.log

[Install]
WantedBy=multi-user.target
